# 8.x/shared/crontab
* * * * * /bin/su -c "/usr/bin/php ${APP_DIRECTORY}/artisan schedule:run" - laravel > /proc/1/fd/1 2>&1
