# 8.x/shared/crontab
* * * * * su-exec laravel php $APP_DIRECTORY/artisan schedule:run > /proc/1/fd/1 2>&1
