# 8.1/shared/crontab
* * * * * /bin/su -c "/usr/bin/php /var/www/html/artisan schedule:run" - laravel > /proc/1/fd/1 2>&1
